# docker/cron/crontab
# handle one batch transfer every minutes
* * * * * root /usr/local/bin/node /app/index.js --network=barnard --senderIndex=0 --count=50 >> /app/cron.log 2>&1
